# scripts/com.merrysails.weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.merrysails.weekly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/resat/Desktop/merrysails/scripts/weekly-review.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/resat/Desktop/merrysails/data/multi-brand-monitoring/weekly-logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/resat/Desktop/merrysails/data/multi-brand-monitoring/weekly-logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/resat/Desktop/merrysails</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
